# services/restart.service
[Unit]
Description=Restart Option Monitor Service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart option-monitor.service

