# services/restart.service
[Unit]
Description=Restart My App Service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart spy-monitor.service

